# crontab
*/5 * * * * echo '=== run script ===' && python /workspace/main.py
