# crontab
0,30 * * * * echo '=== run script ===' && python /workspace/main.py
